# 3-docker-practice/scripts/crontab
#minute hour    day     month   week    user       command
*/2       *       *       *       *       root       /usr/local/bin/backup
*         *       *       *       *       root       echo "${now}_${MYSQL_DATABASE}" > /backup/test
