# 3-docker-practice/scripts/crontab
#minute hour    day     month   week    command
*/1       *       *       *       *       /usr/local/bin/backup
*         *       *       *       *       root       echo hallo

